# systemd/codex-backend.service
[Unit]
Description=installer-codex backend API
After=network.target

[Service]
Type=simple
EnvironmentFile=/opt/installer-codex/.env
ExecStart=/usr/local/bin/installer-codex-run-backend
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
